# recipes-utils/ultra96-startup-pages/files/ultra96-startup-page.service
[Unit]
Description=Start ultra96 python flask server

[Service]
ExecStart=/usr/bin/ultra96-startup-commands.sh

[Install]
WantedBy=multi-user.target
