# recipes-daemons/ospi-dm/files/ov-start.service
[Unit]
Description=ospi-versal

[Service]
ExecStart=/usr/bin/ospi-versal.sh

[Install]
WantedBy=multi-user.target
