# recipes-utils/system-controller-app/files/system_controller.service
[Unit]
Description=system_controller_app

[Service]
ExecStart=/usr/bin/system_controller.sh

[Install]
WantedBy=multi-user.target
