# recipes-utils/start-jupyternb/start-jupyternb/jupyter-setup.service
[Unit]
Description=jupyter setup scripts

[Service]
ExecStart=/sbin/start-jupyter.sh

[Install]
WantedBy=multi-user.target
