# recipes-utils/scweb/files/scwebrun.service
[Unit]
Description=scweb run 

[Service]
ExecStart=/usr/bin/scwebrun.sh

[Install]
WantedBy=multi-user.target
