# recipes-utils/archconfig/files/archconfig.service
[Unit]
Description=archconfig

[Service]
ExecStart=/usr/bin/archconfig.sh

[Install]
WantedBy=multi-user.target
