# recipes-apps/board-fpga-autoload/files/board-fpga-autoload.service
[Unit]
Description=board-fpga-autoload

[Service]
ExecStart=/usr/bin/board-fpga-autoload.sh

[Install]
WantedBy=multi-user.target
