# recipes-utils/boardframework/files/start_boardframework.service
[Unit]
Description=start_boardframework service

[Service]
ExecStart=/usr/bin/start_boardframework.sh

[Install]
WantedBy=multi-user.target
