# recipes-utils/start-jupyterlab/start-jupyterlab/jupyter-setup.service
[Unit]
Description=jupyter setup scripts

[Service]
ExecStart=/sbin/start-jupyter.sh

[Install]
WantedBy=default.target
