# recipes-utils/cmc/cmc-deploy/cmc-start.service
[Unit]
Description=cmc-deploy

[Service]
ExecStart=/usr/bin/cmc

[Install]
WantedBy=multi-user.target

